# deploy/bot-mont-shk-bot.service
[Unit]
Description=Bot_Mont_SHK Telegram bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/Bot_Mont_SHK
Environment=BOT_CONFIG_FILE=/opt/Bot_Mont_SHK/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/Bot_Mont_SHK/.venv/bin/python main.py
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
