# deploy/bot-mont-shk-web.service
[Unit]
Description=Bot_Mont_SHK web app
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/Bot_Mont_SHK
Environment=BOT_CONFIG_FILE=/opt/Bot_Mont_SHK/.env
ExecStart=/opt/Bot_Mont_SHK/.venv/bin/python main_web.py
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
